# examples/systemd/system/spire.target
[Unit]
Description=SPIRE target allowing to start/stop all spire*@.service instances at once

[Install]
WantedBy=multi-user.target
